# myscripts/etc/local.d/98-IMAmeasure.start
#!/bin/bash
dateNow=$(date +"%F_%H-%M-%S")
cat /sys/kernel/security/ima/ascii_runtime_measurements > /var/log/ima/boot_aggregate-${dateNow}
