# myscripts/etc/local.d/98-tpm2measure.start
#!/bin/bash
dateNow=$(date +"%F_%H-%M-%S")
tpm2_eventlog <( cat /sys/kernel/security/tpm0/binary_bios_measurements ) > /var/log/tpm/tpm2_eventlog-${dateNow}

