# myscripts/etc/local.d/99-boot-complete.start
#!/bin/sh
echo "<1>!!! Boot Complete !!!" > /dev/kmsg
BD="/boot/bootcomplete/"
dateNow=$(date +"%F_%H-%M-%S")
uname -a > ${BD}/BootOK-${dateNow}
