# etc/local.d/99-boot-complete.start
#!/bin/sh
dateNow=$(date +"%F_%H-%M-%S")
echo "<1>!!! Boot Complete !!! $dateNow " > /dev/kmsg
